# novncservice/DEBIAN/postrm
#!/bin/sh
systemctl stop novnc-x11vnc.service
systemctl stop novnc.service
systemctl disable novnc-x11vnc.service
systemctl disable novnc.service

