# novncservice/DEBIAN/postinst
#!/bin/sh
systemctl enable novnc-x11vnc.service
systemctl enable novnc.service

systemctl start novnc-x11vnc.service
systemctl start novnc.service

